# ops/systemd/notify@.service
[Unit]
Description=Failure notification for %i

[Service]
Type=oneshot
# Load NTFY_URL from the same secrets file as the other units
EnvironmentFile=-/etc/zakasnyava-li/secrets.env
ExecStart=/usr/bin/bash -c \
    'curl -fsS --max-time 15 \
    -H "Title: Pipeline failure — %i" \
    -d "zakasnyava-li nightly pipeline FAILED: %i" \
    "${NTFY_URL}"'
